# rdhcpd.service
[Unit]
Description=rDHCP - High-Performance DHCP Server
Documentation=https://github.com/ZerosAndOnesLLC/rDHCP
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/bin/rdhcpd /etc/rdhcpd/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
WatchdogSec=30

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/rdhcpd
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes

# Allow binding to privileged ports
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_RAW

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
